# backend/invoice-backend.service
[Unit]
Description=Invoice Assistant Backend Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/invoice-assistant-v2/backend
Environment="PATH=/var/www/invoice-assistant-v2/backend/venv/bin"
ExecStart=/var/www/invoice-assistant-v2/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8090 --workers 2
Restart=always
RestartSec=10

# 日志
StandardOutput=append:/var/log/invoice-backend/access.log
StandardError=append:/var/log/invoice-backend/error.log

[Install]
WantedBy=multi-user.target